# mcp/mcp-server.service
[Unit]
Description=MCP Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=$REPO_DIR
ExecStart=/usr/bin/python3 $REPO_DIR/src/mcp_server_on_raspi/server.py
Restart=always
RestartSec=10
Environment=PATH=/usr/bin:/usr/local/bin
Environment=PYTHONPATH=$REPO_DIR

[Install]
WantedBy=multi-user.target