# mcp/mcp-server.service
[Unit]
Description=MCP Server for Raspberry Pi
After=network.target

[Service]
User=pi
WorkingDirectory=/home/pi/Ultimate-Raspberry-Pi-Setup/mcp
ExecStart=/usr/bin/python3 /home/pi/Ultimate-Raspberry-Pi-Setup/mcp/src/mcp_server_on_raspi/server.py
Restart=on-failure
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=mcp-server

[Install]
WantedBy=multi-user.target